# prod/launch/supervisor.d/backend.ini
[program:fund_be]
command=gunicorn manage:app -w 1 -k gevent --timeout 1200  -b 0.0.0.0:5000 --log-level=info --reload
autorestart=True
autostart=True
environment=
    PYTHONIOENCODING=utf8,
    PYTHONUNBUFFERED=1,
    FLASK_ENV='development'
directory=/home/fund
redirect_stderr = true
stdout_logfile= /tmp/prod_fund.log
stdout_logfile_maxbytes = 20MB